# deploy/lock28-ops-assistant.service
[Unit]
Description=Lock28 Ops Assistant
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/lock28-assistant
EnvironmentFile=/opt/lock28-assistant/.env
ExecStart=/usr/bin/bash /opt/lock28-assistant/start_assistant.sh private
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
